# system/profile.d/go.sh
#!/bin/sh

# Go path
export GOPATH="${XDG_CONFIG_HOME}/go"