# system/profile.d/go.sh
#!/bin/sh

# Go path
export GOPATH="${HOME}/.config/go"