# system/profile.d/gtk.sh
#!/bin/sh

# Gtk 2.0 path
export GTK2_RC_FILES="${HOME}/.config/gtk-2.0/.gtkrc"