# system/profile.d/jre.sh
#!/bin/sh

# Fix problem with Java applications
export _JAVA_AWT_WM_NONREPARENTING=1