# system/profile.d/xinit.sh
#!/bin/sh

# Xinitrc path
export XINITRC="${XDG_CONFIG_HOME}/X11/.xinitrc"