# system/profile.d/xinit.sh
#!/bin/sh

# Xinitrc path
export XINITRC="${HOME}/.config/X11/.xinitrc"